# tags/start/mixter/util/mixter.crontab
#min	hr	mday	month	wday	user		command
0	5	*	*	*	postgres	vacuum-analyze-db.job
18	6	*	*	*	ispivey		db-backup.job
